# raspberry/appletv-monitor/appletv-monitor.service
[Unit]
Description=Apple TV Monitor for IPS Ambilight Control
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/appletv-monitor
Environment=HOME=/home/pi
ExecStart=/home/pi/pyatv/bin/python /home/pi/appletv-monitor/app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
