# raspberry/appletv-monitor.service
[Unit]
Description=Apple TV monitor for IPS-Ambilight-Control
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
Environment=HOME=/home/pi
WorkingDirectory=/home/pi
ExecStart=/home/pi/pyatv/bin/python /opt/appletv-monitor/appletv-monitor.py --address 192.168.178.145 --listen 0.0.0.0 --port 8091 --atvscript /home/pi/pyatv/bin/atvscript
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi

[Install]
WantedBy=multi-user.target
